# service_files/nagios-nest.service
[Unit]
Description=Flask app to convert coralogix alerts to nagios api
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/nagios-alerts-connector
EnvironmentFile=/home/ubuntu/nagios-alerts-connector/nagios-nest.config
ExecStart=/usr/bin/python3 nagios-alert-listener.py
Restart=always

[Install]
WantedBy=multi-user.target